# crontab
* * * * *  python /app/hello.py >> /crontest.txt
# new line